# usr/lib/systemd/profile.d/70-systemd-shell-extra.sh
version https://git-lfs.github.com/spec/v1
oid sha256:91aebbcc8fe1b2c5e2a658ad67cd91e2eb82e4d9b63c50af6461e9c5c17ba95a
size 1067
