# usr/share/factory/etc/profile.d/locale.sh
version https://git-lfs.github.com/spec/v1
oid sha256:153d848ed51f2774e5a1578ea08e0c8586ecc63f7562697e035b84247edb2f82
size 778
